# quadlets/pathfinder-db.container
[Unit]
Description=Pathfinder PostgreSQL database

[Container]
ContainerName=pathfinder-db
Image=docker.io/postgres:16-alpine
Network=pathfinder.network
Volume=pathfinder-postgres-data.volume:/var/lib/postgresql/data
Environment=POSTGRES_USER=postgres
Environment=POSTGRES_PASSWORD=postgres
Environment=POSTGRES_DB=pathfinder
HealthCmd=pg_isready -U postgres
HealthInterval=5s
HealthTimeout=5s
HealthRetries=5

[Service]
Restart=always

[Install]
WantedBy=default.target
